# systemd/acme-tiny@.timer
[Unit]
Description=Monthly ACME certificate renewal

[Timer]
OnCalendar=07:31:30

[Install]
WantedBy=timers.target

