# systemd/acme-tiny@.timer
[Unit]
Description=Monthly ACME certificate renewal

[Timer]
OnCalendar=monthly

[Install]
WantedBy=timers.target

